# gitplop/supervisord/ii.conf
[program:ii]
command=/usr/local/ii/ii.sh
autostart=true
autorestart=true
redirect_stderr=true
user=git
group=git
directory=/tmp
startsecs=60
umask=002
